# windsurf-update.timer
[Unit]
Description=Check and Update Windsurf Editor Daily

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
